# 99-local.rules
ACTION=="add", RUN+="/usr/local/mount-device.sh"
